# systemd/trevcan-socketcan.service.example
[Unit]
Description=TREVCAN SocketCAN Bootstrap
Wants=systemd-udev-settle.service
After=systemd-udev-settle.service
Before=trevcan-explorer.service

[Service]
Type=oneshot
WorkingDirectory=__WORKDIR__
Environment=TREVCAN_SOCKETCAN_BITRATE=__SOCKETCAN_BITRATE__
ExecStart=__PYTHON__ systemd/socketcan_manager.py --bitrate __SOCKETCAN_BITRATE__
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target