# systemd/trevcan-explorer.service.example
[Unit]
Description=TREVCAN Explorer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__GROUP__
WorkingDirectory=__WORKDIR__
Environment=HOME=__HOME__
Environment=PATH=__PATH__
Environment=PYTHONUNBUFFERED=1
ExecStart=__PYTHON__ start.py --service
Restart=on-failure
RestartSec=5
KillMode=control-group
TimeoutStopSec=20
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target